# recipes-wigwag/maestro/maestro/maestro-watcher.service
[Unit]
Description=Maestro restarter

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart maestro.service

[Install]
WantedBy=network.target