# recipes-wigwag/maestro/maestro/maestro-watcher.path
[Unit]
Description=Monitor the changes to identity.json file and restart maestro

[Path]
PathChanged=EDGE_DATA/edge_gw_config/identity.json
Unit=maestro-watcher.service

[Install]
WantedBy=network.target