# recipes-connectivity/coredns/files/coredns.service
[Unit]
Description=coredns
Wants=kubelet.service
After=kubelet.service
Requires=network-online.target
After=network-online.target

[Service]
Type=simple
Restart=no
ExecStart=EDGE_BIN/launch-coredns.sh

[Install]
WantedBy=
