# recipes-connectivity/coredns/files/coredns-resolv.path
[Unit]
Description=Monitor the changes to resolv.conf file and rebuild coredns resolv.conf file

[Path]
PathChanged=/etc/reslov.conf
PathChanged=/run/systemd/resolve/resolv.conf
PathChanged=/etc/resolv-conf.NetworkManager
Unit=coredns-resolv-watcher.service

[Install]
WantedBy=network.target
