# recipes-containers/kubelet/files/kubelet-watcher.service
[Unit]
Description=kubelet restarter

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart kubelet.service

[Install]
WantedBy=network.target