# recipes-containers/kubelet/files/kubelet.service
[Unit]
Description=Kubelet
Wants=edge-proxy.service
After=edge-proxy.service

[Service]
Restart=always
RestartSec=5
ExecCondition=sh -c '[ -f EDGE_DATA/edge_gw_config/identity.json ]'
ExecStart=EDGE_BIN/launch-kubelet.sh

[Install]
WantedBy=network.target
