# recipes-connectivity/pe-terminal/files/pe-terminal-watcher.service
[Unit]
Description=Relay-term restarter

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart pelion-relay-term.service

[Install]
WantedBy=network.target