# recipes-connectivity/pe-terminal/files/pe-terminal-watcher.path
[Unit]
Description=Monitor the changes to identity.json file and restart relay-term

[Path]
PathChanged=EDGE_ETC/pe-terminal/config.json
Unit=pe-terminal-watcher.service

[Install]
WantedBy=network.target