# recipes-wigwag/mbed-edge-core/files/edge-core-byoc.service
[Unit]
Description=Edge Core
After=systemd-networkd-wait-online.service maestro.service

[Service]
Restart=always
RestartSec=5s
ExecStartPre=mkdir -p /userdata/mbed
ExecStartPre=chown developer -R /userdata/mbed
ExecStart=/wigwag/mbed/edge-core  --http-port 9101 --cbor-conf /userdata/mbed/mbed_cloud_dev_credentials.cbor

[Install]
RequiredBy=network.target
