# recipes-wigwag/relay-term/files/relay-term-watcher.path
[Unit]
Description=Monitor the changes to identity.json file and restart relay-term

[Path]
PathChanged=/wigwag/wigwag-core-modules/relay-term/config/config.json
Unit=relay-term-watcher.service

[Install]
WantedBy=network.target